# charms/trusty/cf-cloud-controller/files/upstart/cf-nats.conf
description "Cloud Foundry NATS"
author "Alexander Prismakov<prismakov@gmail.com>"
start on starting cf-cloudcontroller or runlevel [2345]
stop on runlevel [!2345]
expect daemon
#apparmor load <profile-path>
setuid vcap
setgid vcap
respawn
normal exit 0
chdir /var/lib/cloudfoundry/cfcloudcontroller
exec bundle exec nats-server -c jobs/config/nats.yml -d
